# apps/web/install/linux/nanumpay.service
[Unit]
Description=AgentTree SEA service
After=network-online.target mongod.service mongodb.service
Wants=network-online.target

[Service]
Type=simple
User=nanumpay
Group=nanumpay
WorkingDirectory=/opt/nanumpay
EnvironmentFile=-/etc/default/nanumpay
ExecStart=/opt/nanumpay/dist/nanumpay
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
